# HDServer/installation/hivedrive.service
[Unit]
Description=HiveDrive
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=adam
ExecStart=/bin/bash /usr/share/hivedrive/updateHiveDrive.sh

[Install]
WantedBy=multi-user.target
